# systemd/.config/systemd/user/desktop-user.target
[Unit]
Description=User desktop running
BindsTo=graphical-session.target
