# systemd/.config/systemd/user/update-frece-torr.timer
[Unit]
Description=Updates torr index

[Timer]
OnBootSec=2min
OnUnitActiveSec=60min
AccuracySec=20min

[Install]
WantedBy=default.target
