# systemd/.config/systemd/user/update-locate.timer
[Unit]
Description=Updates locate database
After=default.target

[Timer]
OnBootSec=1min
OnUnitActiveSec=1h
AccuracySec=20min

[Install]
WantedBy=default.target
