# systemd/.config/systemd/user/update-frece-dir.service
[Unit]
Description=Updates directory index

[Service]
Type=oneshot
ExecStart=/home/mulhaq/code/frece/examples/dir_update_fd.sh
CPUSchedulingPolicy=idle
IOSchedulingClass=idle
IOSchedulingPriority=7
Nice=19
