# systemd/.config/systemd/user/xsession.target
[Unit]
Description=X session running
BindsTo=graphical-session.target
