# systemd/.config/systemd/user/yams.service
[Unit]
Description=yams
After=mopidy.service

[Service]
Type=simple
ExecStart=/usr/bin/yams -N

[Install]
WantedBy=default.target
